# Computer Engeneering Bachelor/3_year/Amministrazione di Sistemi/LAS/LAS_Prove_Esame/prova_esame_2021_09_06/roles/Systemd_unit/files/count.service
[Unit]
Description=Start /home/count.sh on boot

[Service]
Type=simple
ExecStart=/home/count.sh
Restart=on-failure

[Install]
WantedBy=basic.target