# Computer Engeneering/3_year/Amministrazione di Sistemi/LAS/LAS_Prove_Esame/prova_esame_2021_09_06/roles/Systemd_unit/files/watchdog.service
[Unit]
Description=Start /urs/bin/watchdog.sh on boot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/watchdog.sh
Restart=always

[Install]
WantedBy=multi-user.target